# ru/2025/04/d5729fa5-d717-56c9-af96-45530459d0f8/clamav-freshclam.service
[Unit]
Description=ClamAV virus database updater
Documentation=https://docs.clamav.net/
# If user wants it run from cron, don't start the daemon.
ConditionPathExists=!/etc/cron.d/clamav-freshclam
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/freshclam -d --foreground=true

[Install]
WantedBy=multi-user.target
