# ru/2025/04/96bdcb5c-a58b-58ae-9e6b-536b49bf1c51/fangfrisch.timer
[Unit]
Description=Automatically refresh the fangfrisch signatures
Documentation=https://rseichter.github.io/fangfrisch/
Requires=fangfrisch.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=2m
Unit=fangfrisch.service

[Install]
WantedBy=timers.target
